# supervisord.conf
[supervisord]
nodaemon=true

[program:apachex]
command=/bin/sh -c "exec /etc/init.d/apache2 start"