# supervisord.conf
[supervisord]
nodaemon=true

[program:apachex]
command=/bin/sh -c "exec /etc/init.d/apache2 start"

[program:chmodx]
command=/bin/sh -c "exec chmod -R 777 /var/www/html"